# Documents/BashScripts/etcCreator/etc/udev/rules.d/00-usb-keyboards.rules
ACTION=="add", ATTRS{idVendor}=="feed", ATTRS{idProduct}=="1307", RUN+="/usr/local/bin/x-session-settings-background", OWNER="rmw"