# Documents/BashScripts/etcCreator/etc/udev/rules.d/drawing-tablet.rules
ACTION=="add", ATTRS{idVendor}=="256c", ATTRS{idProduct}=="006d", RUN+="/usr/local/bin/x-session-settings-background", OWNER="rmw"