# Documents/BashScripts/etcCreator/etc/systemd/system/socialMediaBlocker.service
[Unit]
Description=Social media blocker

[Service]
Type=simple
ExecStart=/usr/bin/bash Documents/BashScripts/only_social_media_in_the_weekend.sh
