# Documents/BashScripts/etcCreator/etc/systemd/system/socialMediaBlocker.timer
[Unit]
Description=Do something daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
