# Documents/BashScripts/etcCreator/etc/systemd/system/user-firefox.slice
[Unit]
Description=Firefox Slice
Before=slices.target

[Slice]
CPUQuota=25%
